# amon-agent
#!/usr/bin/env sh
### BEGIN INIT INFO
# Provides:          amon-agent
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts the Amon agent
# Description:       Amon agent - collects system and process information.

AGENTPATH='/usr/bin/amon-agent.py'
AGENTUSER="amonagent"
PIDPATH="/var/run/amonagent/"


[ -f $AGENTPATH ] || echo "$AGENTPATH not found"

action=$1

case $action in
    start)
        if [ ! -d $PIDPATH ]; then
            mkdir -p $PIDPATH
            chown amonagent:amonagent $PIDPATH
        fi
        
        su $AGENTUSER -c "python $AGENTPATH stop"
        su $AGENTUSER -c "python $AGENTPATH start"
        ;;

    stop)
        su $AGENTUSER -c "python $AGENTPATH stop"
        exit $?
        ;;

    restart)
        $0 stop
        $0 start
        exit $?
        ;;

    status)
        su $AGENTUSER -c "python $AGENTPATH status"
        exit $?
        ;;

    install)
        su -c "python $AGENTPATH install $2"
        exit $?
        ;;

    test_collectors)
        su $AGENTUSER -c "python $AGENTPATH test_collectors"
        ;;

    test_plugins)
        su $AGENTUSER -c "python $AGENTPATH test_plugins"
        exit $?
        ;;

    *)
        echo "Usage: $0 {start|stop|restart|test_plugins|status|test_collectors|install}"
        exit 2
        ;;
esac
