# packaging/debian/postrm
update-rc.d amon-agent remove -f
deluser amonagent
groupdel amonagent
rm -rf /var/log/amonagent
rm -rf /etc/amonagent
pip uninstall amonagent

exit 0
