# packaging/debian/postinst
#! /bin/bash

set -e

case "$1" in
    configure)
        easy_install amonagent
        update-rc.d amon-agent defaults
        adduser --system amonagent --disabled-login --shell /bin/sh --no-create-home --quiet
        chown root:root /etc/init.d/amon-agent
        chown -R amonagent:root /etc/amonagent
        chown -R amonagent:root /var/log/amonagent

        set +e

        sudo /etc/init.d/amon-agent restart
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
    ;;
esac

#DEBHELPER#

exit 0
