# packaging/debian/postrm
#! /bin/bash
set -e

if [ "$1" = purge ]; then
    update-rc.d amon-agent remove
    deluser amonagent
fi

exit 0
