# packaging/rpm/postinst
# start
easy_install amonagent
chown -R amon-agent /etc/amonagent/
chown -R amon-agent /var/log/amonagent/
chown -R amon-agent /var/log/amonagent/amonagent.log
chkconfig --add amon-agent

sudo /etc/init.d/amon-agent restart

exit 0
